# branches/groupoffice-3.5-old/debian-groupoffice-mailserver/etc/cron.d/groupoffice-mailserver
1 3 * * * root php /usr/share/groupoffice-mailserver/cache_usage.php