# trunk/debian-groupoffice-servermanager/etc/cron.d/groupoffice-servermanager
1 4 * * * root php /usr/share/groupoffice/modules/servermanager/cache_report.php