# branches/groupoffice-4.1/debian-groupoffice-com/etc/cron.d/groupoffice-cron
5 * * * * www-data php /usr/share/groupoffice/groupofficecli.php -r=core/cron -c=/etc/groupoffice/config.php > /dev/null
