# branches/groupoffice-com-3.02/debian-groupoffice-mailserver/etc/cron.d/groupoffice-mailserver
1 3 * * * root php /usr/share/groupoffice/modules/postfixadmin/cache_usage.php