# debian-groupoffice-com/etc/cron.d/groupoffice-cron
5 * * * * root php /usr/share/groupoffice/groupofficecli.php -r=core/cron -c=/etc/groupoffice/config.php > /dev/null