# branches/groupoffice-4.0/debian-groupoffice-mailserver/etc/cron.d/groupoffice-mailserver
#Updates usage information in the servermanager.
1 3 * * * root /usr/bin/php /usr/share/groupoffice-mailserver/cache_usage.php
