# trunk/debian-groupoffice-documents/etc/cron.d/groupoffice-documents
0 1 * * * root /usr/bin/php /usr/share/groupoffice/groupofficecli.php -c=/etc/groupoffice/config.php -r=filesearch/filesearch/sync > /dev/null
